# daemon/pocketmeter-ble-daemon.service
[Unit]
Description=PocketMeter BLE Daemon
Requires=bluetooth.target
After=bluetooth.target

[Service]
Type=simple
ExecStart=DAEMON_PATH
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
